# templates/mail/mail.service
[Unit]
Description=Fetch mail

[Service]
Type=oneshot
ExecStart=/home/{{username}}/.bin/mail.sh
